# docker/etc/periodic/daily/certbot
#!/bin/sh
certbot --webroot -w /var/www/foxcaves/html/static --post-hook 's6-svc -wr -r /etc/s6/nginx' renew
