# docker/etc/periodic/daily/certbot
#!/bin/sh
exec certbot --nginx --nginx-server-root /usr/local/openresty/nginx/conf renew 
