# docker/etc/periodic/15min/expiry_cleanup
#!/bin/sh

exec curl --unix-socket /run/nginx-lua-api.sock -s 'http://main.foxcaves/api/v1/system/expire' -XPOST
