# etc/s6/redis/run
#!/bin/sh

cd /var/lib/redis
exec /usr/bin/redis-server
