# etc/s6/nginx/run
#!/bin/sh

mkdir -p /tmp/fc_thumbs /tmp/fc_files
chown www-data:www-data /tmp/fc_thumbs /tmp/fc_files /opt/foxcaves_storage

exec /usr/local/openresty/nginx/sbin/nginx -g 'daemon off; user www-data;'
