# etc/periodic/15min/expiry_cleanup
#!/bin/sh

exec curl --resolve "foxcaves.foxden.network:443:127.0.0.1" -s "https://foxcaves.foxden.network/api/v1/system/expire" -XPOST
