# etc/periodic/15min/expiry_cleanup
#!/bin/sh

exec curl --resolve "foxcav.es:443:127.0.0.1" -s "https://foxcav.es/api/v1/system/expire" -XPOST
