# etc/periodic/15min/expiry_cleanup
#!/bin/sh

exec curl --resolve "caves.fox.ax:443:127.0.0.1" -s "https://caves.fox.ax/api/v1/system/expire" -XPOST
