# etc/periodic/15min/expiry_cleanup
#!/bin/sh

exec curl --resolve "fcv.doridian.net:443:127.0.0.1" -s "https://fcv.doridian.net/api/v1/system/expire" -XPOST
