# etc/periodic/15min/expiry_cleanup
#!/bin/sh
curl -s "https://${HOSTNAME}/api/v1/system/expire" -XPOST
