# docker/etc/s6/nginx/run
#!/bin/sh
set -e

export CAPTCHA_FONT=/usr/share/fonts/opensans/OpenSans-Regular.ttf

mkdir -p /var/www/foxcaves/tmp /var/www/foxcaves/storage /etc/letsencrypt/storage
chown foxcaves:foxcaves /var/www/foxcaves/tmp /var/www/foxcaves/storage /etc/letsencrypt/storage

openssl req -x509 -newkey rsa:2048 -keyout /etc/letsencrypt/snakeoil.key -out /etc/letsencrypt/snakeoil.crt -sha256 -days 3650 -nodes -subj '/CN=snakeoil'

if [ ! -f /etc/letsencrypt/account.key ]; then
    openssl genpkey -algorithm RSA -pkeyopt rsa_keygen_bits:4096 -out /etc/letsencrypt/account.key
fi

luajit /var/www/foxcaves/lua/nginx_configure.lua

rm -f /run/nginx-lua-api.sock

exec /usr/local/openresty/bin/openresty -c /usr/local/openresty/nginx/conf/custom.conf
